# sys_config/lib/systemd/system/radio.service
[Unit]
Description=KRadio
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/pi/raspiradio/prog/radio.py