# sys_config/etc/init.d/radio
#! /bin/bash
### BEGIN INIT INFO
# Provides: radio
# Required-Start: $all
# Required-Stop:
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: radio program
# Description:
### END INIT INFO
case "$1" in
start)
        echo -n "Starting radio..."
        start-stop-daemon --start --quiet --chuid kradio --group kradio --make-pidfile --pidfile /var/run/kradio.pid --exec /home/kradio/raspiradio/prog/radio.py
        echo "radio started"
        ;;
    stop)
        echo -n "Stopping radio..."
        start-stop-daemon --stop --quiet --pidfile /var/run/kradio.pid
        killall mpd
        echo "stopped"
        ;;
    *)
        echo "Use /etc/init.d/radio {start|stop}"
        exit 1
        ;;
esac
exit 0

